# bot/systemd/ecosystem-watch.timer
[Unit]
Description=Daily ZAO ecosystem watch
[Timer]
OnCalendar=*-*-* 13:00:00 UTC
Persistent=true
[Install]
WantedBy=timers.target
