# bot/systemd/fleet-heartbeat.service
[Unit]
Description=Fleet heartbeat poster -> thezao.xyz bots board
[Service]
Type=oneshot
ExecStart=%h/fleet-heartbeat/heartbeat.sh
