# bot/systemd/fleet-claude-auth.timer
[Unit]
Description=claude auth check every 10 min
[Timer]
OnBootSec=120
OnUnitActiveSec=10min
[Install]
WantedBy=timers.target
