# bot/systemd/fleet-claude-auth.service
[Unit]
Description=Check claude CLI auth + alert Zaal on failure
[Service]
Type=oneshot
ExecStart=%h/fleet-heartbeat/claude-auth-check.sh
